# data/deploy/loki.service
[Unit]
Description=Loki log aggregation
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/loki -config.file=/etc/loki/config.yml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
